# scripts/bridge-watchdog.service
[Unit]
Description=Bridge Watchdog — health check and self-healing
After=robothor-bridge.service

[Service]
Type=oneshot
ExecStart=/home/philip/robothor/brain/memory_system/venv/bin/python /home/philip/robothor/brain/scripts/bridge_watchdog.py
User=philip
WorkingDirectory=/home/philip/robothor/brain
